# autothread.service
[Unit]
Description=Neon Auto-Threading Daemon
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/clawd/skills/telegram-userbot
ExecStart=/home/ubuntu/clawd/.venv/bin/python /home/ubuntu/clawd/skills/telegram-userbot/autothread_daemon.py --daemon --interval 30
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
